# packaging/systemd/kiri-sync.timer
[Unit]
Description=Daily push of kiri notes to GitHub

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
